# brave.cron
0 0 * * * python /home/leai/barve/src_one/PersistentOffLine/UpdateData.py 
