# brave.cron
0 0 * * 1 python /home/leai/barve/src_one/PersistentOffLine/UpdateData.py 
